# installers/configfiles/skytracer.socket
[Unit]
Description=skytracer socket
[Socket]
ListenStream=/run/skytracer.sock
[Install]
WantedBy=sockets.target
